# solr.conf
[program:start-solr]
user = solr
directory = /opt/solr
command = /bin/bash bin/start-solr-node.sh
autorestart = true
